# dist/init/upstart/mangband.conf
# mangband - launches the mangband server
#

description     "Launches MAngband server"

env USER=mangband
env GROUP=games
env CMD=/usr/games/mangband
env LOGFILE=/var/log/mangband/mangband.log
env PIDFILE=/var/run/mangband/mangband.pid

start on runlevel [2345]
stop on runlevel [06]
respawn
respawn limit 10 5
kill signal INT
kill timeout 10

# on upstart < 1.10, the server will respawn on stop, sorry
# normal exit 255 # doesn't help :(

exec start-stop-daemon --start --user $USER --group $GROUP \
    --background --make-pidfile --pidfile $PIDFILE \
    --exec /bin/bash -- -c "exec $CMD > $LOGFILE 2>&1"

pre-stop script
    start-stop-daemon --signal INT -R INT/1/INT/1/INT/1/INT/1/INT/1 --stop --exec "$CMD" --pidfile $PIDFILE
end script
